# templates/default/omero-web-init.erb
#!/bin/bash
# omero         Start/stop script for omero.
#               Harri Jäälinoja <harri.jaalinoja@helsinki.fi> 
#
# chkconfig:    345 97 5
# description:  OMERO web server
  
OMERO_USER=<%= @omero_user %>
OMERO_HOME=<%= @omero_home %>
  
. /etc/init.d/functions
RETVAL=0
  
case "$1" in
  start | restart )
    [ "$1" = start ] && _P=start || _P=restart
  
    # start (or restart) OMERO server
    runuser -s /bin/bash - $OMERO_USER -c \
      "ulimit -c 81920; $OMERO_HOME/bin/omero web $_P"
  
    RETVAL=$?; echo; [ $RETVAL -eq 0 ] && touch /var/lock/subsys/`basename $0`
    exit $RETVAL ;;
  
  stop )
    runuser -s /bin/bash - $OMERO_USER -c "$Env $OMERO_HOME/bin/omero web stop"
    RETVAL=$?; echo; [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/`basename $0`
    exit $RETVAL ;;

  status )
    runuser -s /bin/bash - $OMERO_USER -c "$Env $OMERO_HOME/bin/omero web status"
    RETVAL=$?; echo
    exit $RETVAL ;;

  * )
    echo "Usage: $0 {start|stop|status|restart}"
    exit 2 ;;
esac
